# tree/etc/cron.d/tt-rss-update-feeds
#Cron job to update TT-rss feeds hourly.
0 * * * * www-data /usr/bin/php /var/www/tt-rss/update.php --feeds --quiet
30 * * * * www-data /usr/bin/php /var/www/tt-rss/update.php --feeds --quiet
