# tree/etc/cron.weekly/logwatch
#Weekly log analysis
59 23 * * 5	root	logwatch --range 'between -7 days and today' --filename /var/log/last-logwatch

